# flood-sensor.service
[Unit]
Description=Flood Sensor Monitoring Service
After=network.target
Wants=network.target

[Service]
Type=forking
User=pi
Group=pi
WorkingDirectory=~/Desktop/Flood-Sensor
ExecStart=~/Desktop/Flood-Sensor/run.sh start
ExecStop=~/Desktop/Flood-Sensor/run.sh stop
ExecReload=~/Desktop/Flood-Sensor/run.sh restart
PIDFile=~/Desktop/Flood-Sensor/flood_sensor.pid
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
